# clawgate.service
[Unit]
Description=ClawGate – LLM Routing Proxy for OpenClaw
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=clawgate
Group=clawgate
WorkingDirectory=/opt/clawgate
EnvironmentFile=/opt/clawgate/.env
Environment=CLAWGATE_DB_PATH=/var/lib/clawgate/clawgate.db
ExecStart=/opt/clawgate/venv/bin/uvicorn \
    clawgate.main:app \
    --host 127.0.0.1 \
    --port 8090 \
    --workers 1 \
    --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/clawgate
PrivateTmp=true

[Install]
WantedBy=multi-user.target
