# docs/examples/faigate-auto-update.cron
# Example cron entry for controlled fusionAIze Gate auto-update checks.
# Runs every day at 03:30 and only applies updates that the runtime marks as eligible.
30 3 * * * faigate /usr/local/bin/faigate-auto-update --apply >>/var/log/faigate-auto-update.log 2>&1
