# docs/examples/faigate-auto-update.service
[Unit]
Description=Apply eligible fusionAIze Gate updates
Documentation=https://github.com/fusionAIze/faigate
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=faigate
Group=faigate
WorkingDirectory=/opt/faigate
EnvironmentFile=/opt/faigate/.env
ExecStart=/usr/local/bin/faigate-auto-update --apply
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/faigate /opt/faigate
PrivateTmp=true
